# sflow-prom.container
[Container]
ContainerName=sflow-prom
Image=docker.io/sflow/prometheus
Pod=monitoring.Pod
